# ops/jody-contact-api.service
[Unit]
Description=Jody Portfolio Contact API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/jody/contact-api
EnvironmentFile=/etc/jody/contact-api.env
ExecStart=/usr/bin/node /var/www/jody/contact-api/dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
